# .config/systemd/user/mpd-mpris.service.d/10-bind-to-mpd.conf
[Unit]
PartOf=mpd.service

[Install]
#Clearing WantedBy before adding a new service to it
WantedBy=
WantedBy=mpd.service
